# autostart/init_synaptiks_config.desktop
[Desktop Entry]
Type=Application
Exec=synaptikscfg init
Name=Initialize touchpad configuration
OnlyShowIn=KDE;

X-KDE-autostart-phase=1
